# elbot-update.service.template
[Unit]
Description=Auto-update ELBOT
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=__USER__
WorkingDirectory=__WORKDIR__
Environment=PATH=__WORKDIR__/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=__WORKDIR__/deploy.sh
StandardOutput=append:/var/log/elbot/auto-update.log
StandardError=append:/var/log/elbot/auto-update.log
SuccessExitStatus=0
RemainAfterExit=yes
