# elbot-update.service
[Unit]
Description=Auto-update ELBOT

[Service]
Type=oneshot
ExecStart=/home/alex/ELBOT/deploy.sh
