# elbot.service.template
[Unit]
Description=ELBOT Discord Bot with Web Interface
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__WORKDIR__
ExecStart=__WORKDIR__/venv/bin/python __WORKDIR__/main.py
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3
CPUQuota=40%
MemoryMax=384M
EnvironmentFile=__WORKDIR__/.env
Environment=ELBOT_WEB_PORT=8081
Environment=FLASK_ENV=production

[Install]
WantedBy=multi-user.target
